# scraper/vf-poller.service
[Unit]
Description=VF Poller service (Tippmix-v-sport)
After=network.target

[Service]
Type=simple
User=youruser
WorkingDirectory=/path/to/your/project    # <-- cseréld le a projekt gyökérre (ahol a scraper könyvtár van)
ExecStart=/usr/bin/python3 -u -m scraper.cli poll
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target